# docs/service/linux-systemd.service
[Unit]
Description=Codex Feishu Bridge
After=network-online.target

[Service]
Type=simple
Environment=FEISHU_APP_SECRET=change-me
ExecStart=/usr/local/bin/codex-feishu-bridge serve --config /home/example/.codex-feishu-bridge/config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
