# docs/service/macos-launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>info.fofa.codex-feishu-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>codex-feishu-bridge</string>
    <string>serve</string>
    <string>--config</string>
    <string>/Users/example/.codex-feishu-bridge/config.toml</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/example/.codex-feishu-bridge/bridge.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/example/.codex-feishu-bridge/bridge.err.log</string>
</dict>
</plist>
